# internxt/rootfs/etc/services.d/internxt/run
#!/usr/bin/env bashio

# Read configuration
INXT_USER=$(bashio::config 'inxt_user')
INXT_PASSWORD=$(bashio::config 'inxt_password')
INXT_TWOFACTORCODE=$(bashio::config 'inxt_twofactorcode')
INXT_OTPTOKEN=$(bashio::config 'inxt_otptoken')
WEBDAV_PORT=$(bashio::config 'webdav_port')
WEBDAV_PROTOCOL=$(bashio::config 'webdav_protocol')

bashio::log.info "Starting Internxt WebDAV server..."

# Check required configuration
if bashio::config.is_empty 'inxt_user'; then
    bashio::exit.nok "Internxt username/email is required!"
fi

if bashio::config.is_empty 'inxt_password'; then
    bashio::exit.nok "Internxt password is required!"
fi

# Set environment variables for the Internxt CLI
export INXT_USER="${INXT_USER}"
export INXT_PASSWORD="${INXT_PASSWORD}"
export WEBDAV_PORT="${WEBDAV_PORT}"
export WEBDAV_PROTOCOL="${WEBDAV_PROTOCOL}"

# Handle 2FA configuration
if ! bashio::config.is_empty 'inxt_otptoken'; then
    export INXT_OTPTOKEN="${INXT_OTPTOKEN}"
    bashio::log.info "Using OTP token for automatic 2FA"
elif ! bashio::config.is_empty 'inxt_twofactorcode'; then
    export INXT_TWOFACTORCODE="${INXT_TWOFACTORCODE}"
    bashio::log.info "Using provided 2FA code"
fi

bashio::log.info "Configuration validated"
bashio::log.info "Internxt User: ${INXT_USER}"
bashio::log.info "WebDAV Protocol: ${WEBDAV_PROTOCOL}"
bashio::log.info "WebDAV Port: ${WEBDAV_PORT}"

# Login to Internxt (non-interactive using env vars)
bashio::log.info "Logging into Internxt..."
internxt login

# Start the WebDAV server
bashio::log.info "Enabling Internxt WebDAV server..."
exec internxt webdav enable
